# userpatches/overlay/rootfs/etc/profile.d/99-fix-serial-term.sh
# fix terminal for UART user
if [ "$TERM" == "unknown" ]; then
    TERM="xterm"
fi

export TERM
